# templates/main.timer.j2
[Unit]
After=multi-user.target

[Timer]
OnCalendar={{ systemd_timer_frequency }}
Persistent=yes

[Install]
WantedBy=timers.target
