# templates/main.timer.j2
[Unit]
After=multi-user.target

[Timer]
OnCalendar={{ systemd_timer_frequency }}
RandomizedDelaySec={{ systemd_timer_random_delay_sec }}
Persistent=yes

[Install]
WantedBy=timers.target
